# deploy/realtor-agent-sync.service.example
[Unit]
Description=Realtor Agent scheduled sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=/path/to/Realtor_Agent
EnvironmentFile=/path/to/Realtor_Agent/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/path/to/Realtor_Agent/.venv/bin/realtor-agent --scheduled-sync --db-path /path/to/Realtor_Agent/data/realtor_agent.db
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
